# 1-Homer_Nginx/Conf/crontab
@reboot sleep 10 && /mnt/Cloud/Data/iptables_custom.sh